# deploy/macos/com.gohttpauto.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gohttpauto.worker</string>
  <key>ProgramArguments</key>
  <array>
    <string>__APP_DIR__/server/gohttpauto</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__APP_DIR__/server</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
    <!-- DEBUG: visible Chrome — remove GFX_VISIBLE / GFX_KEEP_OPEN when done checking -->
    <key>GFX_VISIBLE</key>
    <string>1</string>
    <key>GFX_KEEP_OPEN</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/gohttpauto/worker.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/gohttpauto/worker.err.log</string>
</dict>
</plist>
