# deploy/macos/com.gohttpauto.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gohttpauto.tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/ssh</string>
    <string>-i</string>
    <string>__HOME__/.ssh/vps_tunnel</string>
    <string>-N</string>
    <string>-o</string>
    <string>ServerAliveInterval=30</string>
    <string>-o</string>
    <string>ServerAliveCountMax=3</string>
    <string>-o</string>
    <string>ExitOnForwardFailure=yes</string>
    <string>-o</string>
    <string>BatchMode=yes</string>
    <string>-L</string>
    <string>3307:127.0.0.1:3306</string>
    <string>root@74.208.99.161</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/gohttpauto/tunnel.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/gohttpauto/tunnel.err.log</string>
</dict>
</plist>
